# Sample systemd/cloudflare-dns-refresh.timer
[Unit]
Description=Daily DNS Refresh

[Timer]
OnCalendar=*-*-* 00:00:00
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target